# raspberry-pi/dustcheck.service
[Unit]
Description=DustCheck Air Quality Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/dustcheck
ExecStart=/home/pi/dustcheck/dust_env/bin/python -u monitor.py
Restart=always
RestartSec=10
StandardOutput=append:/home/pi/dustcheck/monitor.log
StandardError=append:/home/pi/dustcheck/monitor.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
